# deployment/wix-printer.service
[Unit]
Description=Wix Printer Service (FastAPI)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=wix-printer
Group=wix-printer
WorkingDirectory=/opt/wix-printer-service
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/wix-printer-service/venv/bin/uvicorn wix_printer_service.api.main:app --host 0.0.0.0 --port 8000 --log-level info
Restart=always
RestartSec=3
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target
